# content/service/4/run
#!/bin/sh

EXEC=$(echo $RANDOM | md5sum | head -c 6; echo)
ln -s /mnt/data/config/rclone.conf /tmp/4${EXEC}.conf
mv /workdir/4* /workdir/4${EXEC}

# Run rclone
htpasswd -csb /workdir/htpasswd ${GLOBAL_USER} ${GLOBAL_PASSWORD}
exec 2>&1
exec /workdir/4${EXEC} rcd --rc-serve \
--rc-web-gui-no-open-browser \
--rc-htpasswd=/workdir/htpasswd \
--rc-addr=localhost:61802 \
--buffer-size 4M \
--transfers 3 \
--config /tmp/4${EXEC}.conf \
--rc-allow-origin=* \
--use-mmap
