# content/service/caddy/run
#!/bin/sh

# Configure AriaNG & Rclone WebUI
sed -i 's|6800|'"${CADDY_WEB_PORT}"'|g' /workdir/ariang/js/aria-ng-3fcc0271ed.min.js
sed -i 's|rpcInterface:"jsonrpc"|rpcInterface:"'"${GLOBAL_PORTAL_PATH/\//}"'/jsonrpc"|g' /workdir/ariang/js/aria-ng-3fcc0271ed.min.js
sed -i 's|ipAddress:a|ipAddress:`${window.location.protocol}//${window.location.hostname}:${window.location.port}'"${GLOBAL_PORTAL_PATH}"'/rclonerc`|g' /workdir/rcloneweb/build/static/js/3.90421092.chunk.js

# Configure homer
cp /workdir/homer_conf/* /workdir/homer/assets/tools/
cp /workdir/homer_conf/homer_${GLOBAL_LANGUAGE}.yml /workdir/homer/assets/config.yml
sed -i "s|GLOBAL_PORTAL_PATH|${GLOBAL_PORTAL_PATH}|g" /workdir/homer/assets/config.yml

if [ ! -f "/workdir/homer_conf/pyload.png" ]; then
    LINE=$(sed -n -e '/pyLoad/=' /workdir/homer/assets/config.yml)
    LINE2=$((LINE + 5))
    sed -i "${LINE},${LINE2}d" /workdir/homer/assets/config.yml
fi

if [ "${OLIVETIN}" != "Enable" ]; then
    LINE=$(sed -n -e '/\"OliveTin\"/=' /workdir/homer/assets/config.yml)
    LINE2=$((LINE + 5))
    sed -i "${LINE},${LINE2}d" /workdir/homer/assets/config.yml
fi

# Configure Caddyfile
mkdir -p /mnt/data/config/caddy 2>/dev/null
if [ ! -f "/mnt/data/config/caddy/Caddyfile" ]; then
    cp /workdir/Caddyfile /mnt/data/config/caddy/
fi

cp -f /mnt/data/config/caddy/Caddyfile /tmp/Caddyfile
HASH="$(caddy hash-password --plaintext ${GLOBAL_PASSWORD})"
sed -i "s|HASH|${HASH}|g" /tmp/Caddyfile

# Run caddy
exec 2>&1
exec caddy run --config /tmp/Caddyfile --adapter caddyfile
